# simp/submodules/nofever/nofever/config/rc.local
#!/bin/sh -e

sleep 20
sudo -H -u nofever0003 sudo /home/nofever0003/NoFever/nofever/boot.sh &

exit 0


